# netgen-pro-dpdk.service
[Unit]
Description=NetGen Pro - DPDK Edition High-Performance Traffic Generator
After=network.target network-online.target
Wants=network-online.target
Documentation=https://github.com/your-repo/netgen-pro-dpdk

[Service]
Type=simple
User=root
WorkingDirectory=/opt/netgen-dpdk/web
Environment="PATH=/opt/netgen-dpdk/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
Environment="VIRTUAL_ENV=/opt/netgen-dpdk/venv"

# Ensure hugepages are available
ExecStartPre=/bin/bash -c 'if [ ! -d /mnt/huge ]; then mkdir -p /mnt/huge; mount -t hugetlbfs nodev /mnt/huge; fi'
ExecStartPre=/bin/bash -c 'echo 1024 > /sys/kernel/mm/hugepages/hugepages-2048kB/nr_hugepages || true'

# Load kernel modules
ExecStartPre=/sbin/modprobe -a uio vfio-pci || true

# Start the web control server
ExecStart=/opt/netgen-dpdk/venv/bin/python /opt/netgen-dpdk/web/dpdk_control_server.py

# Restart behavior
Restart=always
RestartSec=10

# Cleanup on stop
ExecStop=/bin/bash -c 'pkill -TERM -f dpdk_engine || true'
ExecStopPost=/bin/bash -c 'rm -f /tmp/dpdk_engine_control.sock || true'

# Security
NoNewPrivileges=false
PrivateTmp=false

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=netgen-pro-dpdk

[Install]
WantedBy=multi-user.target
